# debian/postinst
#!/bin/sh
set -e

DATA_DIR="/var/lib/keko-ts3bot"
CONFIG_DIR="/etc/keko-ts3bot"
CONFIG_FILE="$CONFIG_DIR/keko-ts3bot.yaml"

case "$1" in
    configure)
        echo ""
        echo "============================================"
        echo "  Kellerkompanie TeamSpeak 3 Bot"
        echo "============================================"
        echo ""

        # Create teamspeak3bot system user if it doesn't exist
        if ! id -u teamspeak3bot >/dev/null 2>&1; then
            echo "Creating system user 'teamspeak3bot'..."
            useradd --system --no-create-home --shell /usr/sbin/nologin teamspeak3bot
        fi

        # Create data directory
        mkdir -p "$DATA_DIR"
        chown teamspeak3bot:teamspeak3bot "$DATA_DIR"
        chmod 750 "$DATA_DIR"

        # Restrict config file - contains database and query credentials
        chown root:teamspeak3bot "$CONFIG_FILE"
        chmod 640 "$CONFIG_FILE"

        # Reload systemd daemon
        systemctl daemon-reload

        # Enable but do not start - config must be edited first
        systemctl enable keko-ts3bot

        echo "Service enabled (not started - configure first)."
        echo ""
        echo "Configuration: $CONFIG_FILE"
        echo ""
        echo "  1. Edit the configuration file:"
        echo "     sudo nano $CONFIG_FILE"
        echo "  2. Start the service:"
        echo "     sudo systemctl start keko-ts3bot"
        echo "  3. View logs:"
        echo "     journalctl -u keko-ts3bot -f"
        echo ""
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
